# .vscode-server/data/User/History/4f3775a/Zq80.conf
[program:trading-bot]
directory=/root/bot
command=/root/bot/myenv/bin/python3 main.py  # Full path to Python executable in the virtual environment
user=root
autostart=true
autorestart=true  # Consider setting this to true if you want the bot to restart automatically on failure
stderr_logfile=/root/trading-bot.log
stdout_logfile=/root/trading-bot.log
